# deploy/velox-edge.service
[Unit]
Description=Velox Edge Trading Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/velox-edge
ExecStart=/opt/velox-edge/.venv/bin/python run.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
